# hardware/ops/twinr-peer-ai-camera-proxy.service
[Unit]
Description=Twinr peer AI camera observation proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/twinr-peer-ai-camera/repo
ExecStart=/opt/twinr-peer-ai-camera/venv/bin/python /opt/twinr-peer-ai-camera/repo/hardware/ops/peer_ai_camera_observation_proxy.py --repo-root /opt/twinr-peer-ai-camera/repo --env-file /opt/twinr-peer-ai-camera/repo/.proxy.env --host 10.42.0.2 --port 8767
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
