# hardware/ops/twinr-remote-memory-watchdog.service
[Unit]
Description=Twinr Remote Memory Watchdog
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=0
ConditionPathExists=/proc/device-tree/model

[Service]
Type=simple
WorkingDirectory=/twinr
Environment=TWINR_PYTHON=/twinr/.venv/bin/python
ExecCondition=/usr/bin/env bash -lc 'grep -qa "Raspberry Pi" /proc/device-tree/model'
ExecStart=/usr/bin/env bash -lc 'cd /twinr && PYTHONPATH=src "${TWINR_PYTHON}" -m twinr --env-file .env --watch-remote-memory'
Restart=always
RestartSec=2
RestartPreventExitStatus=75
TimeoutStopSec=20
KillMode=control-group
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
