# hardware/ops/twinr-remote-memory-watchdog.service
[Unit]
Description=Twinr Remote Memory Watchdog
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=thh
Group=thh
WorkingDirectory=/twinr
Environment=TWINR_PYTHON=/twinr/.venv/bin/python
ExecStart=/usr/bin/env bash -lc 'cd /twinr && PYTHONPATH=src "${TWINR_PYTHON}" -m twinr --env-file .env --watch-remote-memory'
Restart=always
RestartSec=2
TimeoutStopSec=20
KillMode=control-group
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
