# hardware/ops/twinr-whatsapp-channel.service
[Unit]
Description=Twinr WhatsApp Channel
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=thh
WorkingDirectory=/twinr
Environment=PYTHONPATH=src
ExecStart=/twinr/.venv/bin/python -m twinr --env-file .env --run-whatsapp-channel
Restart=always
RestartSec=2
KillSignal=SIGINT
SuccessExitStatus=130 143
TimeoutStopSec=30
KillMode=control-group
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
